# test/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
priority=10
directory=/
command=/usr/sbin/sshd -D
user=marc
autostart=true
autorestart=true
stopsignal=QUIT


[program:xvfb]
priority=10
directory=/
command=/usr/bin/Xvfb :1 -screen 0 1024x768x16
user=marc
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err


[program:x11vnc]
priority=10
directory=/
command=x11vnc -display :1 -xkb
user=marc
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err
