# test/supervisord.conf
##################################################################################################################################################################
# Supervisord #
###############
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/fd/1
logfile_maxbytes=0


##################################################################################################################################################################
# TigerVNC #
############
[program:x11vnc]
priority=20
directory=/
command=x11vnc -display :1 -xkb -auth /var/run/lightdm/root/:0 -forever
user=marc
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/x11vnc.log
redirect_stderr=true

##################################################################################################################################################################
# NoVNC (Alternative: websockify) #
###################################
